# install/service-templates/vibe.service
[Unit]
Description=Vibe Daemon (dotvibe.dev)
Documentation=https://dotvibe.dev
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/dotvibe/vibe daemon
Restart=always
RestartSec=10
User=root
Group=root
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=/root
WorkingDirectory=/var/lib/vibe
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vibe-daemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vibe

[Install]
WantedBy=multi-user.target